# conf/examples/rc.local
# for System Activity Grapher
#
su -l sag -c 'PATH-TO-CMD/bin/addgap'
#
touch /tmp/ntpctl.out
chown sag:sag /tmp/ntpctl.out
chmod 0640 /tmp/ntpctl.out
